# mopidy/user-mopidy.service
[Unit]
Description=mopidy
After=network.target sound.target

[Service]
ExecStart=/usr/bin/mopidy

[Install]
WantedBy=default.target
