# sabnzbd/sabnzbd.user.service
[Unit]
Description=SABnzbd binary newsreader

[Service]
ExecStart=/opt/sabnzbd/SABnzbd.py

[Install]
WantedBy=default.target
